# unitree/go1/deploy/go1-beep-adapter.service
[Unit]
Description=Go1 Beep Adapter (head speaker beep/volume) — serves :18082 /v1/beep/actions
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=unitree
# 绑 192.168.123.13:18082；开机 IP 未就绪时由 Restart 反复重试。beep 专属，与 speaker 的
# go1-head-adapter(:18080) 分离——只做 beep/set_volume/get_volume。
ExecStart=/usr/bin/python3 /home/unitree/beep_adapter.py --config /home/unitree/go1-beep-adapter.json
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
